# infra/systemd/door-sync.service
[Unit]
Description=Doorboard durable sync service
Requires=doorboard-storage-ready.service
Wants=door-media.service
After=doorboard-storage-ready.service door-media.service

[Service]
Type=simple
User=doorboard
Group=doorboard
WorkingDirectory=/opt/doorboard
EnvironmentFile=/etc/doorboard/doorboard.env
ExecStart=/opt/doorboard/.venv/bin/python -m door_sync
Restart=on-failure
RestartSec=2s
TimeoutStopSec=20s
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/mnt/ssd/doorboard -/mnt/nas/doorboard

[Install]
WantedBy=multi-user.target
