# infra/systemd/door-media.service
[Unit]
Description=Doorboard media service
Requires=doorboard-storage-ready.service
Wants=mediamtx.service
After=doorboard-storage-ready.service mediamtx.service

[Service]
Type=simple
User=doorboard
Group=doorboard
WorkingDirectory=/opt/doorboard
EnvironmentFile=/etc/doorboard/doorboard.env
ExecStart=/opt/doorboard/.venv/bin/python -m door_media
Restart=on-failure
RestartSec=2s
TimeoutStopSec=15s
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/mnt/ssd/doorboard

[Install]
WantedBy=multi-user.target
